# cron.daily/pkgfile-update
#!/bin/sh

pkgfile --update && echo "[$(date)] : pkgfile update successful" >> /tmp/cron.pkgfile-update.log

